# servers/supervisor/celerybeat.blade.ini
; Python Celery Beat scheduler config

[program:{{ $worker->name }}]
directory={{ $project->projectDir }}

user={{ $project->serverUsername }}
numprocs=1
stdout_logfile=/var/log/celery/{{ $worker->name }}.log
stderr_logfile=/var/log/celery/{{ $worker->name }}.log
autostart=true
autorestart=true
startsecs=10

command={{ $worker->command() }}

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs= {{ $worker->stopSeconds ?? '600' }}

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Set Celery priority higher than default (999)
; so, if a broker program (redis/rabbitmq) is supervised, it will start first.
priority=999
